# supervisor-worker.conf
[program:exam-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/clubit-exam/htdocs/exam.clubit.id/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/home/clubit-exam/htdocs/exam.clubit.id/storage/logs/worker.log
stopwaitsecs=3600
